# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:simple-host-checker]
command=/usr/bin/pwsh /app/simple-host-checker.ps1